# app/refactored_app/deploy/us-tickers-update.timer
[Unit]
Description=Daily update for US tickers CSV (NYSE+NASDAQ)

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
