# app/refactored_app/deploy/us-tickers-update.service
[Unit]
Description=Update US tickers CSV (NASDAQ + major U.S. exchanges) for the Streamlit app
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=pi
ExecStart=/home/pi/venv/bin/python /home/pi/refactored_app_current/deploy/build_us_ticker_dict.py
